# contrib/init/alphacashd.conf
description "Alphacash Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env ALPHACASHD_BIN="/usr/bin/alphacashd"
env ALPHACASHD_USER="alphacash"
env ALPHACASHD_GROUP="alphacash"
env ALPHACASHD_PIDDIR="/var/run/alphacashd"
# upstart can't handle variables constructed with other variables
env ALPHACASHD_PIDFILE="/var/run/alphacashd/alphacashd.pid"
env ALPHACASHD_CONFIGFILE="/etc/alphacash/alphacash.conf"
env ALPHACASHD_DATADIR="/var/lib/alphacashd"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # alphacashd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$ALPHACASHD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run alphacashd."
        echo "The setting must appear in $ALPHACASHD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Alphacash Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$ALPHACASHD_PIDDIR"
    chmod 0755 "$ALPHACASHD_PIDDIR"
    chown $ALPHACASHD_USER:$ALPHACASHD_GROUP "$ALPHACASHD_PIDDIR"
    chown $ALPHACASHD_USER:$ALPHACASHD_GROUP "$ALPHACASHD_CONFIGFILE"
    chmod 0660 "$ALPHACASHD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$ALPHACASHD_PIDFILE" \
    --chuid $ALPHACASHD_USER:$ALPHACASHD_GROUP \
    --exec "$ALPHACASHD_BIN" \
    -- \
    -pid="$ALPHACASHD_PIDFILE" \
    -conf="$ALPHACASHD_CONFIGFILE" \
    -datadir="$ALPHACASHD_DATADIR" \
    -disablewallet \
    -daemon

